# tasks/linux/systemd-timer/monithor.timer
[Unit]
Description=Monithor's daily timer

[Timer]
Unit=monithor.service
OnCalendar=daily

[Install]
WantedBy=multi-user.target
